# ofborg/deploy/tickborg@.service
[Unit]
Description=Tickborg %i service
After=network-online.target rabbitmq-server.service
Wants=network-online.target
PartOf=tickborg.target

[Service]
Type=simple
User=tickborg
Group=tickborg
WorkingDirectory=/var/lib/tickborg
ExecStart=/var/lib/tickborg/bin/%i /etc/tickborg/config.json

Environment=RUST_LOG=tickborg=info
Environment=RUST_BACKTRACE=1
Environment=HOME=/var/lib/tickborg

# Sandbox
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/tickborg /var/log/tickborg

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tickborg-%i

[Install]
WantedBy=tickborg.target
